# supervisord.d/quick_startup.ini
[program:quick_startup]
command=./home/sc/桌面/syscuts/src/boot_startup.sh
numprocs=1
autostart=true
user=root
redirect_stderr=true
stdout_logfile=/home/sc/桌面/syscuts/log/boot_startup.log
